# build/airootfs/etc/NetworkManager/dispatcher.d/99-macspoof
#!/bin/bash
# Randomize MAC address on every network interface bring-up
[ "$2" = "up" ] || exit 0
ip link set "$1" down
macchanger -r "$1" 2>/dev/null
ip link set "$1" up
